# BashScripts/crontab
*/10 * * * * /home/ubuntu/track/pushFails.sh